# renew.crontab
0 0,12 * * * python -c 'import random; import time; tim.sleep(random.random() * 3600)' && certbot renew
